# deploy/systemd/asset-lifecycle.service
[Unit]
Description=IT Asset Lifecycle Manager
After=network.target

[Service]
Type=simple
User=asset-lifecycle
Group=asset-lifecycle
WorkingDirectory=/opt/asset-lifecycle-manager
EnvironmentFile=/etc/asset-lifecycle/asset-lifecycle.env
ExecStart=/opt/asset-lifecycle-manager/.venv/bin/python -m uvicorn main:app --app-dir /opt/asset-lifecycle-manager/backend --host 127.0.0.1 --port 8000 --workers 1
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
